# overlay/etc/sv/.httpd/run
#!/bin/bash

exec python3 -m http.server --directory /opt/www/ 80 &
