# overlay/etc/sv/mailhog/run
#!/bin/bash

exec mailhog -api-bind-addr 0.0.0.0:8025 -ui-bind-addr 0.0.0.0:8025 -smtp-bind-addr 0.0.0.0:1025 > /dev/null