# overlay/etc/sv/samba/run
#!/bin/bash

exec nmbd -D && smbd -FS --no-process-group
